# bot.service
[Unit]
Description=Discord Bot
After=network.target

[Service]
WorkingDirectory=
Type=simple
User=
ExecStart=/bin/bash ./start.sh
Restart=always

[Install]
WantedBy=multi-user.target
