# bot.service
[Unit]
Description=Discord Bot
After=network.target

[Service]
WorkingDirectory=
Type=simple
User=
ExecStart=node index.js
Restart=always

[Install]
WantedBy=multi-user.target